# Ressourcen/plant_watering.service
[Unit]
Description=Pflanzen-Bewässerungssystem
After=network.target

[Service]
ExecStart=/path/to/plant_watering.sh
Restart=always
User=pi
WorkingDirectory=/home/pi

[Install]
WantedBy=multi-user.target